# files/etc/supervisord.d/3-bootstrap.conf
[program:bootstrap]
command=sh /%(program_name)s.sh
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
priority=3
user=cocorico