# files/etc/supervisord.d/cocorico.conf
[program:cocorico]
command=sh /%(program_name)s.sh
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
user=%(program_name)s
priority=2