# files/etc/supervisord.d/cocorico.conf
[program:cocorico]
command = /bin/sh /cocorico.sh
redirect_stderr=true
user = cocorico
