# files/etc/supervisord.d/cocorico.conf
[program:cocorico]
command = /bin/sh /cocorico.sh
stdout_logfile = tmp/cocorico.log
user = cocorico
