# files/etc/supervisord.d/mongo.conf
[program:mongo]
command = sh /mongo.sh
stdout_logfile = tmp/mongo.log
