# files/etc/supervisord.d/mongo.conf
[program:mongo]
command = sh /mongo.sh
redirect_stderr=true
