# files/etc/supervisord.d/mysql.conf
[program:mysql]
command = sh /mysql.sh
redirect_stderr=true
