# files/etc/supervisord.d/mysql.conf
[program:mysql]
command = sh /mysql.sh
stdout_logfile = tmp/mysql.log
